# printforge/scripts/ustreamer.service
[Unit]
Description=ustreamer MJPEG camera streamer
After=network.target

[Service]
ExecStart=/opt/printforge/ustreamer \
    --device /dev/video0 \
    --host 0.0.0.0 \
    --port 8080 \
    --resolution 640x480 \
    --desired-fps 30 \
    --persistent \
    --quality 80
Restart=always
RestartSec=3
User=david1534
Group=video

[Install]
WantedBy=multi-user.target
